# setup/vpn_port_rpcserver.conf
[program:vpn_port_rpcserver]
user = nobody
command = /vpn-ip-responder/vpn_port_rpcserver -l 0.0.0.0 -p 8108 -c /config
